# deploy/timers/example.timer
[Unit]
Description=Example Service

[Timer]
OnCalendar=Tue 11:20:00
Persistent=true

[Install]
WantedBy=timers.target
